# heliox-mon.service
[Unit]
Description=Heliox Monitor - 服务器监控系统
Documentation=https://github.com/hh/heliox-mon
After=network.target

[Service]
Type=simple
User=root
Group=root

# 工作目录
WorkingDirectory=/opt/heliox-mon

# 环境变量配置文件
EnvironmentFile=-/opt/heliox-mon/.env

# 启动命令
ExecStart=/usr/local/bin/heliox-mon

# 重启策略
Restart=always
RestartSec=5

# 资源限制
LimitNOFILE=65535
LimitNPROC=65535

# 安全加固
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/heliox-mon /var/lib/heliox-mon

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=heliox-mon

[Install]
WantedBy=multi-user.target
